# transcriber/com.localloop.transcriber.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.localloop.transcriber</string>

    <key>ProgramArguments</key>
    <array>
        <string>/path/to/your/venv/bin/python3</string>
        <string>/path/to/local-loop/transcriber/transcribe.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/path/to/local-loop/transcriber</string>

    <!-- Run daily at 2:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/localloop-transcriber.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/localloop-transcriber.error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
